# forge/specs/ermete-recovery/ermete-recovery-1.0.0/systemd/ermete-recovery.target
[Unit]
Description=Ermete OS Emergency Recovery Mode Target
Documentation=https://github.com/patapem/ermete-forge
Requires=ermete-recovery.service
AllowIsolate=yes
