# system_files/etc/systemd/system/btrfs-maintenance.timer
[Unit]
Description=Monthly BTRFS Maintenance Timer

[Timer]
OnCalendar=monthly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
