# build_files/usr/lib/systemd/user/niri-session.target
[Unit]
Description=Niri Wayland Session
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target

[Target]
