# forge/specs/ermete-ui-agent/ermete-ui-agent-1.0.0/ermete-ui-agent.service
[Unit]
Description=Ermete Generative UI Agent
PartOf=graphical-session.target

[Service]
ExecStart=/usr/libexec/ermete-ui-agent/agent.py
Restart=always

[Install]
WantedBy=default.target
