# system_files/usr/lib/systemd/user/ermete-ide-bootstrap.service
[Unit]
Description=Ermete OS: Bootstrap Nix IDE Stack (gcc, make, lazygit)
After=network-online.target niri.service
Wants=network-online.target
ConditionPathExists=!%h/.local/state/ermete-ide-bootstrapped

[Service]
Type=oneshot
ExecStart=/usr/libexec/ermete-ide-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
