# system_files/usr/lib/systemd/user/swaync.service
[Unit]
Description=SwayNC (Wayland Notification Daemon)
PartOf=graphical-session.target
After=graphical-session-pre.target

[Service]
Type=simple
ExecStart=/usr/bin/swaync
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
