# specs/ermete-system-config/SOURCES/usr/lib/systemd/system/ermete-timewarp.timer
[Unit]
Description=Timer for Ermete Btrfs Time-Warp Snapshots

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
