# forge/specs/ermete-system-services/SOURCES/usr/lib/systemd/user/niri-session.target
[Unit]
Description=Ermete OS Niri Desktop Session
Documentation=https://github.com/patapem/ermete-forge
BindsTo=graphical-session.target
Wants=graphical-session.target xdg-desktop-autostart.target
After=graphical-session.target xdg-desktop-autostart.target
