# system_files/usr/lib/systemd/system/greetd.service.d/10-ermete.conf
[Service]
ExecStart=
ExecStart=/usr/bin/greetd --config /usr/lib/greetd/config.toml
