# system_files/usr/lib/systemd/system/greetd.service.d/10-wait-drm.conf
[Unit]
Wants=dev-dri-renderD128.device
After=dev-dri-renderD128.device
