# forge/specs/ermete-base-config/SOURCES/usr/lib/systemd/system/nvidia-powerd.service.d/laptop-only.conf
[Unit]
# Impedisce il crash del servizio su PC Desktop non supportati
ConditionChassis=|laptop
ConditionChassis=|convertible
